# scripts/crontab

@reboot /bin/bash /usr/bin/brscan-skey
